# debian/dnssec-root-key.cron.weekly
#!/bin/sh
unbound-anchor -a /etc/dns/root.key \
               -c /usr/share/dnssec-root-key/icannbundle.pem
